# supervisord.conf
[supervisord]
nodaemon=true

[program:couchdb]
command=couchdb

[program:lamp]
command=/run.sh
